# files/system/etc/systemd/system/bootc-fetch-apply-updates.service.d/10-no-apply.conf
[Service]
ExecStart=
ExecStart=/usr/bin/bootc update --quiet
